# s6-overlay/s6-rc.d/avahi-log/run
#!/usr/bin/env bash
# Logger for the avahi service. s6 pipes avahi's stdout into here; s6-log writes
# rotated, timestamped files to /mt2mqtt-run/log/avahi/current.
#   n5      keep # rotated files
#   s1000000 rotate at ~1 MB
set -euo pipefail

mkdir -p /mt2mqtt-run/log            # s6-log creates the per-service dir itself
exec s6-log -b n5 s1000000 /mt2mqtt-run/log/avahi
