# s6-overlay/s6-rc.d/bluetoothd-log/run
#!/usr/bin/env bash
# Logger for the bluetoothd service. s6 pipes bluetoothd's stdout into here;
# s6-log writes rotated, timestamped files to /var/log/bluetoothd/current. All
# service logs live under /var/log so a single bind mount exposes them on the host.
#   n5       keep 5 rotated files
#   s1000000 rotate at ~1 MB
set -euo pipefail

exec s6-log -b n5 s1000000 /var/log/bluetoothd    # s6-log creates the leaf dir itself
