# s6-overlay/s6-rc.d/dbus-log/run
#!/usr/bin/env bash
# Logger for the dbus service. s6 pipes dbus's stdout into here; s6-log writes
# rotated, timestamped files to /mt2mqtt-run/logs/dbus/current. All service logs
# live under /mt2mqtt-run/logs (in the bind-mounted run dir) so they show on the host.
#   n5       keep 5 rotated files
#   s1000000 rotate at ~1 MB
set -euo pipefail

mkdir -p /mt2mqtt-run/logs                  # s6-log makes only the leaf, not the parent
exec s6-log -b n5 s1000000 /mt2mqtt-run/logs/dbus
