# s6-overlay/s6-rc.d/matter-server/run
#!/usr/bin/env bash
set -euo pipefail

exec 2>&1                         # fold stderr into stdout so matter-server-log captures it

args=(
    --storage-path /mt2mqtt-run/matter-server-data
    --paa-root-cert-dir /src/connectedhomeip/credentials/production/paa-root-certs
    --listen-address 127.0.0.1
    --port 5580
)

# Enable BLE commissioning from whatever BT path is available, in priority order:
# 1. The HOST's BlueZ, over a D-Bus proxy at $HOST_BUS. xdg-dbus-proxy should be
#    configured to expose org.bluez and nothing else from the host. The services
#    in the container can use BT, but don't have access to anything else host-side.
#    Chiptool handles BT over dbus, so this is enough for commissioning a new device
#    using the host's BT stack.
# 2. A BT adapter directly visible in THIS netns (host-net fallback): use the
#    container's own bluetoothd on the container's own system bus.
# 3. Neither: run as an operational bridge, no BLE and no commissioning.

HOST_BUS=/mt2mqtt-run/bluez-proxy.sock
if bt-host-check.sh "$HOST_BUS"; then
    echo "[matter-server] host BlueZ present on $HOST_BUS: BLE via host bluetoothd"
    export DBUS_SYSTEM_BUS_ADDRESS="unix:path=$HOST_BUS"
    args+=(--bluetooth-adapter 0)
elif compgen -G '/sys/class/bluetooth/hci*' >/dev/null; then
    echo "[matter-server] local BT adapter present: BLE via container BlueZ"
    args+=(--bluetooth-adapter 0)
else
    echo "[matter-server] no usable BlueZ (host bus or local adapter); no BLE commissioning"
fi

exec matter-server "${args[@]}"
