# s6-overlay/s6-rc.d/matter-server-log/run
#!/usr/bin/env bash
# Logger for the matter-server service. s6 pipes matter-server's stdout into here;
# s6-log writes rotated, timestamped files to /mt2mqtt-run/logs/matter-server/current.
# All service logs live under /mt2mqtt-run/logs (in the bind-mounted run dir) so
# they show on the host.
#   n5        keep 5 rotated files
#   s10000000 rotate at ~10 MB
set -euo pipefail

mkdir -p /mt2mqtt-run/logs                             # s6-log makes only the leaf, not the parent
exec s6-log -b n5 s10000000 /mt2mqtt-run/logs/matter-server
