# s6-overlay/s6-rc.d/otbr-agent-log/run
#!/usr/bin/env bash
# Logger for the otbr-agent service. s6 pipes otbr-agent's stdout into here;
# s6-log writes rotated, timestamped files to /mt2mqtt-run/logs/otbr-agent/current.
# All service logs live under /mt2mqtt-run/logs (in the bind-mounted run dir) so
# they show on the host.
#   n5        keep 5 rotated files
#   s10000000 rotate at ~10 MB
set -euo pipefail

LOGDIR=/mt2mqtt-run/logs/otbr-agent
mkdir -p "$LOGDIR"                     # make the dir ourselves so the symlink has a home
ln -sf current "$LOGDIR/otbr-agent.log"  # stable .log handle to the live log; survives rotation
exec s6-log -b n5 s10000000 "$LOGDIR"
