# s6-overlay/s6-rc.d/bluetoothd/run
#!/usr/bin/env bash
# bluetoothd (BLE stack). This is the container's OWN BlueZ, kept only
# as a fallback. Commisioning may happen through this BT stack, or through the host's
# BT stack by going through D-Bus.

exec 2>&1                         # fold stderr into stdout so bluetoothd-log captures it

# A BT controller shows up as /sys/class/bluetooth/hci* (netns-scoped). None -> no adapter.
if ! compgen -G '/sys/class/bluetooth/hci*' >/dev/null; then
    echo "[bluetoothd] in-container BT disabled (no adapter found). Commisioning needs to happen through host BT."
    exec s6-svc -d /run/service/bluetoothd    # stay down; don't restart-loop
fi

echo "[bluetoothd] BT adapter found; starting bluetoothd..."
exec /usr/libexec/bluetooth/bluetoothd -n -d
